# obmenik.conf
[program:obmenbot]
directory=/home/ubuntu/ObmenBot
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile=/home/ubuntu/ObmenBot/logfile_err.log
stdout_logfile=/home/ubuntu/ObmenBot/logfile.log